# apps/mac-bridge/launchd/ai.internjobs.mac-bridge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  ai.internjobs.mac-bridge — Mac-side bridge for iMessage traffic.

  DEPENDENCY: BlueBubbles Server must be running locally (127.0.0.1:1234)
  before this LaunchAgent starts. mac-bridge connects to BlueBubbles over
  HTTP + Socket.IO and will retry on backoff if BlueBubbles is down, but
  for clean boots make sure BlueBubbles.app is in your Login Items (or
  managed via the Phase 07b Wave-2 Astrolabe LaunchDaemon).

  The bridge itself runs as a LaunchAgent (per-user) because BlueBubbles
  must run as the logged-in user to access Messages.app — running the
  bridge in the same session avoids cross-user TCC issues.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.internjobs.mac-bridge</string>

  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/node</string>
    <string>--env-file=/Users/raj/internjobs-cms/apps/mac-bridge/.env</string>
    <string>/Users/raj/internjobs-cms/apps/mac-bridge/src/server.mjs</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/raj/internjobs-cms/apps/mac-bridge</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/opt/homebrew/sbin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <dict>
    <key>SuccessfulExit</key>
    <false/>
    <key>Crashed</key>
    <true/>
  </dict>

  <key>ThrottleInterval</key>
  <integer>10</integer>

  <key>StandardOutPath</key>
  <string>/Users/raj/Library/Logs/internjobs-mac-bridge.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/raj/Library/Logs/internjobs-mac-bridge.err.log</string>

  <key>ProcessType</key>
  <string>Interactive</string>
</dict>
</plist>
