# apps/mac-bridge/launchd/ai.internjobs.mac-bridge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.internjobs.mac-bridge</string>

  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/node</string>
    <string>/Users/raj/internjobs-cms/apps/mac-bridge/src/server.mjs</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/raj/internjobs-cms/apps/mac-bridge</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/opt/homebrew/sbin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <dict>
    <key>SuccessfulExit</key>
    <false/>
    <key>Crashed</key>
    <true/>
  </dict>

  <key>ThrottleInterval</key>
  <integer>10</integer>

  <key>StandardOutPath</key>
  <string>/Users/raj/Library/Logs/internjobs-mac-bridge.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/raj/Library/Logs/internjobs-mac-bridge.err.log</string>

  <key>ProcessType</key>
  <string>Interactive</string>
</dict>
</plist>
